# rootfs/etc/s6-overlay/s6-rc.d/php/finish
#!/usr/bin/with-contenv bash
# Actions après l'arrêt de PHP
bashio::log.info "PHP service stopped with exit code $1"