# rootfs/etc/s6-overlay/s6-rc.d/php/finish
#!/usr/bin/with-contenv bash

echo "PHP service stopped with exit code $1"