# rootfs/etc/s6-overlay/s6-rc.d/php/run
#!/usr/bin/with-contenv bashio

bashio::log.info  "Launching PHP..."

exec php -F