# rootfs/etc/services.d/mariadb/run
#!/usr/bin/with-contenv bash

echo "Starting MariaDB..."

mkdir -p /data/mysql

chown -R mysql:mysql /data/mysql

mysqld --datadir=/data/mysql