# rootfs/etc/s6-overlay/s6-rc.d/test/run
#!/usr/bin/with-contenv bash

bashio::log.info  "Creating database..."
echo "coucou"