# rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio

bashio::net.wait_for 9001
bashio::log.info "Starting NGinx..."

exec nginx