# rootfs/etc/s6-overlay/s6-rc.d/mariadb/finish
#!/usr/bin/with-contenv bash
echo "MariaDB service stopped with exit code $1"