# rootfs/etc/s6-overlay/s6-rc.d/mariadb/run
#!/usr/bin/with-contenv bash

bashio::log.info "Starting MariaDB..."

exec mkdir -p /data/mysql

exec chown -R mysql:mysql /data/mysql

exec mysqld --datadir=/data/mysql