# rootfs/etc/s6-overlay/s6-rc.d/mariadb/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting MariaDB..."

mkdir -p /data/mysql /run/mysqld
chown -R mysql:mysql /data/mysql /run/mysqld

exec mysqld --datadir=/data/mysql --user=root